# projects/legacy-archive/_archive/imports/soc-content/detection-rules/sigma/execution/scheduled_task_execution.yml
title: Suspicious Scheduled Task Execution
id: 2d3e4f5a-6b7c-8d9e-0f1a-2b3c4d5e6f7a
status: stable
description: Detects execution of suspicious scheduled tasks
references:
    - https://attack.mitre.org/techniques/T1053/005/
author: HawkinsOps SOC
date: 2025/01/15
modified: 2025/01/15
tags:
    - attack.execution
    - attack.persistence
    - attack.t1053.005
logsource:
    category: process_creation
    product: windows
detection:
    selection:
        ParentImage|endswith: '\taskeng.exe'
    selection_schtasks:
        Image|endswith: '\schtasks.exe'
        CommandLine|contains: '/Create'
    condition: 1 of selection_*
falsepositives:
    - Legitimate scheduled tasks
    - System maintenance
level: low
