# projects/legacy-archive/_archive/imports/soc-content/detection-rules/wazuh/hawkinsops-custom-rules.xml
<!--
  HawkinsOps Custom Wazuh Detection Rules
  Author: HawkinsOps SOC Team
  Date: 2025-01-15
  Description: Production-ready detection rules for SIEM deployment
  Rule ID Range: 100001-100200
-->

<group name="hawkinsops,">

  <!-- CREDENTIAL ACCESS RULES (100001-100020) -->

  <rule id="100001" level="12">
    <if_group>sysmon_event_10</if_group>
    <field name="sysmon.targetImage">\\lsass.exe</field>
    <field name="sysmon.grantedAccess">0x1410|0x1010|0x1438|0x143a|0x1418</field>
    <description>Suspicious LSASS process access detected (possible credential dumping)</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
  </rule>

  <rule id="100002" level="15">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\rundll32.exe</field>
    <field name="sysmon.commandLine">comsvcs.dll</field>
    <field name="sysmon.commandLine">MiniDump</field>
    <description>LSASS memory dump via comsvcs.dll detected</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
  </rule>

  <rule id="100003" level="15">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">sekurlsa::|privilege::debug|lsadump::|invoke-mimikatz</field>
    <description>Mimikatz credential dumping tool detected</description>
    <mitre>
      <id>T1003</id>
    </mitre>
  </rule>

  <rule id="100004" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\reg.exe</field>
    <field name="sysmon.commandLine">save</field>
    <field name="sysmon.commandLine">HKLM\\sam|HKLM\\system|HKLM\\security</field>
    <description>SAM registry hive export detected</description>
    <mitre>
      <id>T1003.002</id>
    </mitre>
  </rule>

  <rule id="100005" level="15">
    <if_sid>60103</if_sid>
    <field name="win.eventdata.properties">1131f6aa-9c07-11d1-f79f-00c04fc2dcd2|1131f6ad-9c07-11d1-f79f-00c04fc2dcd2</field>
    <description>DCSync attack detected - Directory replication abuse</description>
    <mitre>
      <id>T1003.006</id>
    </mitre>
  </rule>

  <rule id="100006" level="10">
    <if_sid>60106</if_sid>
    <field name="win.eventdata.ticketEncryptionType">0x17</field>
    <description>Potential Kerberoasting attack - RC4 ticket request</description>
    <mitre>
      <id>T1558.003</id>
    </mitre>
  </rule>

  <rule id="100007" level="12">
    <if_group>sysmon_process_creation</if_group>
    <match>ntdsutil|ntds.dit</match>
    <description>NTDS.dit credential extraction attempt</description>
    <mitre>
      <id>T1003.003</id>
    </mitre>
  </rule>

  <!-- PERSISTENCE RULES (100021-100040) -->

  <rule id="100021" level="8">
    <if_sid>60102</if_sid>
    <description>New scheduled task created</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
  </rule>

  <rule id="100022" level="10">
    <if_group>sysmon_event_13</if_group>
    <field name="sysmon.targetObject">\\Microsoft\\Windows\\CurrentVersion\\Run</field>
    <description>Registry run key modified for persistence</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <rule id="100023" level="12">
    <if_sid>60612</if_sid>
    <field name="win.system.eventID">7045</field>
    <description>New Windows service installed</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
  </rule>

  <rule id="100024" level="12">
    <if_group>wmi</if_group>
    <match>EventConsumer|EventFilter</match>
    <description>WMI event subscription persistence detected</description>
    <mitre>
      <id>T1546.003</id>
    </mitre>
  </rule>

  <rule id="100025" level="8">
    <if_group>sysmon_event_11</if_group>
    <field name="sysmon.targetFilename">\\Start Menu\\Programs\\Startup\\</field>
    <description>File created in startup folder</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <!-- PRIVILEGE ESCALATION RULES (100041-100060) -->

  <rule id="100041" level="12">
    <if_group>sysmon_event_13</if_group>
    <field name="sysmon.targetObject">ms-settings\\shell\\open\\command</field>
    <description>UAC bypass via fodhelper detected</description>
    <mitre>
      <id>T1548.002</id>
    </mitre>
  </rule>

  <rule id="100042" level="12">
    <if_group>sysmon_event_13</if_group>
    <field name="sysmon.targetObject">mscfile\\shell\\open\\command</field>
    <description>UAC bypass via eventvwr detected</description>
    <mitre>
      <id>T1548.002</id>
    </mitre>
  </rule>

  <rule id="100043" level="10">
    <if_sid>60106</if_sid>
    <field name="win.eventdata.privilegeList">SeDebugPrivilege|SeImpersonatePrivilege</field>
    <description>Suspicious privilege assignment detected</description>
    <mitre>
      <id>T1134</id>
    </mitre>
  </rule>

  <rule id="100044" level="12">
    <if_group>sysmon_event_13</if_group>
    <field name="sysmon.targetObject">AlwaysInstallElevated</field>
    <field name="sysmon.details">0x00000001</field>
    <description>AlwaysInstallElevated registry modification</description>
    <mitre>
      <id>T1548.002</id>
    </mitre>
  </rule>

  <!-- DEFENSE EVASION RULES (100061-100080) -->

  <rule id="100061" level="12">
    <if_sid>60104</if_sid>
    <field name="win.system.eventID">1102</field>
    <description>Windows event log cleared</description>
    <mitre>
      <id>T1070.001</id>
    </mitre>
  </rule>

  <rule id="100062" level="12">
    <if_group>sysmon_event_13</if_group>
    <field name="sysmon.targetObject">Windows Defender\\DisableAntiSpyware|DisableAntiVirus</field>
    <description>Windows Defender disabled via registry</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
  </rule>

  <rule id="100063" level="12">
    <if_group>sysmon_event_8</if_group>
    <description>Process injection via remote thread creation</description>
    <mitre>
      <id>T1055</id>
    </mitre>
  </rule>

  <rule id="100064" level="15">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\svchost.exe|\\lsass.exe|\\csrss.exe</field>
    <not_field name="sysmon.image">^C:\\Windows\\System32\\|^C:\\Windows\\SysWOW64\\</not_field>
    <description>Process masquerading as Windows system process detected</description>
    <mitre>
      <id>T1036</id>
    </mitre>
  </rule>

  <rule id="100065" level="15">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">AmsiUtils|amsiInitFailed|AmsiScanBuffer</field>
    <description>AMSI bypass attempt detected</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
  </rule>

  <rule id="100066" level="10">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">netsh</field>
    <field name="sysmon.commandLine">advfirewall</field>
    <field name="sysmon.commandLine">off</field>
    <description>Windows Firewall disabled</description>
    <mitre>
      <id>T1562.004</id>
    </mitre>
  </rule>

  <rule id="100067" level="15">
    <if_group>sysmon_event_6</if_group>
    <field name="sysmon.signed">false</field>
    <description>Unsigned driver loaded - potential rootkit</description>
    <mitre>
      <id>T1014</id>
    </mitre>
  </rule>

  <!-- LATERAL MOVEMENT RULES (100081-100100) -->

  <rule id="100081" level="8">
    <if_sid>60612</if_sid>
    <field name="win.eventdata.serviceName">PSEXESVC</field>
    <description>PsExec service installation detected</description>
    <mitre>
      <id>T1021.002</id>
    </mitre>
  </rule>

  <rule id="100082" level="6">
    <if_sid>60106</if_sid>
    <field name="win.eventdata.logonType">10</field>
    <description>Remote Desktop logon detected</description>
    <mitre>
      <id>T1021.001</id>
    </mitre>
  </rule>

  <rule id="100083" level="8">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.parentImage">\\wmiprvse.exe</field>
    <description>Process execution via WMI</description>
    <mitre>
      <id>T1047</id>
    </mitre>
  </rule>

  <rule id="100084" level="6">
    <if_sid>60106</if_sid>
    <field name="win.eventdata.shareName">\\C$|\\ADMIN$|\\IPC$</field>
    <description>Windows admin share accessed</description>
    <mitre>
      <id>T1021.002</id>
    </mitre>
  </rule>

  <rule id="100085" level="10">
    <if_sid>60106</if_sid>
    <field name="win.eventdata.logonType">3</field>
    <field name="win.eventdata.logonProcessName">NtLmSsp</field>
    <description>Potential Pass-the-Hash attack</description>
    <mitre>
      <id>T1550.002</id>
    </mitre>
  </rule>

  <!-- EXECUTION RULES (100101-100120) -->

  <rule id="100101" level="10">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">DownloadString|Invoke-Expression|IEX|EncodedCommand</field>
    <description>Suspicious PowerShell command execution</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
  </rule>

  <rule id="100102" level="10">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">-e |-en |-enc |-encoded</field>
    <description>Base64 encoded PowerShell command</description>
    <mitre>
      <id>T1027</id>
    </mitre>
  </rule>

  <rule id="100103" level="10">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\wscript.exe|\\cscript.exe</field>
    <field name="sysmon.commandLine">.vbs|.js|.jse</field>
    <description>Script execution via WScript/CScript</description>
    <mitre>
      <id>T1059.005</id>
    </mitre>
  </rule>

  <rule id="100104" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\mshta.exe</field>
    <field name="sysmon.commandLine">http://|https://|javascript:|vbscript:</field>
    <description>MSHTA abuse detected</description>
    <mitre>
      <id>T1218.005</id>
    </mitre>
  </rule>

  <rule id="100105" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\regsvr32.exe</field>
    <field name="sysmon.commandLine">scrobj.dll|/i:http</field>
    <description>Regsvr32 squiblydoo attack detected</description>
    <mitre>
      <id>T1218.010</id>
    </mitre>
  </rule>

  <rule id="100106" level="8">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.parentImage">\\WINWORD.EXE|\\EXCEL.EXE|\\POWERPNT.EXE</field>
    <field name="sysmon.image">\\cmd.exe|\\powershell.exe|\\wscript.exe</field>
    <description>Office macro spawned suspicious process</description>
    <mitre>
      <id>T1204.002</id>
    </mitre>
  </rule>

  <!-- DISCOVERY RULES (100121-100135) -->

  <rule id="100121" level="4">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\whoami.exe</field>
    <description>User discovery via whoami</description>
    <mitre>
      <id>T1033</id>
    </mitre>
  </rule>

  <rule id="100122" level="4">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\net.exe|\\net1.exe</field>
    <field name="sysmon.commandLine">view|group|localgroup</field>
    <description>Network reconnaissance command</description>
    <mitre>
      <id>T1018</id>
    </mitre>
  </rule>

  <rule id="100123" level="4">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\systeminfo.exe</field>
    <description>System information discovery</description>
    <mitre>
      <id>T1082</id>
    </mitre>
  </rule>

  <rule id="100124" level="6">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">Get-ADUser|Get-ADGroup|Get-ADComputer</field>
    <description>Active Directory enumeration</description>
    <mitre>
      <id>T1087.002</id>
    </mitre>
  </rule>

  <!-- COLLECTION RULES (100136-100150) -->

  <rule id="100136" level="10">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">Get-Clipboard|Windows.Forms.Clipboard</field>
    <description>Clipboard data collection</description>
    <mitre>
      <id>T1115</id>
    </mitre>
  </rule>

  <rule id="100137" level="10">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">CopyFromScreen|System.Drawing.Bitmap</field>
    <description>Screen capture activity detected</description>
    <mitre>
      <id>T1113</id>
    </mitre>
  </rule>

  <rule id="100138" level="12">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">GetAsyncKeyState|SetWindowsHookEx</field>
    <description>Keylogging activity detected</description>
    <mitre>
      <id>T1056.001</id>
    </mitre>
  </rule>

  <rule id="100139" level="8">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">Compress-Archive|7z.exe a|rar.exe a</field>
    <description>Data archiving for collection</description>
    <mitre>
      <id>T1560</id>
    </mitre>
  </rule>

  <!-- EXFILTRATION RULES (100151-100165) -->

  <rule id="100151" level="8">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">dropbox.com|drive.google.com|mega.nz</field>
    <description>Data exfiltration to cloud service</description>
    <mitre>
      <id>T1567</id>
    </mitre>
  </rule>

  <rule id="100152" level="10">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\ftp.exe</field>
    <field name="sysmon.commandLine">put</field>
    <description>FTP upload detected</description>
    <mitre>
      <id>T1048.003</id>
    </mitre>
  </rule>

  <rule id="100153" level="6">
    <if_group>powershell</if_group>
    <field name="win.eventdata.scriptBlockText">Send-MailMessage|-Attachments</field>
    <description>Email-based data exfiltration</description>
    <mitre>
      <id>T1048.003</id>
    </mitre>
  </rule>

  <!-- IMPACT RULES (100166-100185) -->

  <rule id="100166" level="15">
    <if_group>sysmon_event_11</if_group>
    <field name="sysmon.targetFilename">.encrypted|.locked|.crypto|.locky</field>
    <description>Ransomware file encryption detected</description>
    <mitre>
      <id>T1486</id>
    </mitre>
  </rule>

  <rule id="100167" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">vssadmin</field>
    <field name="sysmon.commandLine">delete shadows</field>
    <description>Volume shadow copy deletion</description>
    <mitre>
      <id>T1490</id>
    </mitre>
  </rule>

  <rule id="100168" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\bcdedit.exe</field>
    <field name="sysmon.commandLine">recoveryenabled no</field>
    <description>Boot configuration modified</description>
    <mitre>
      <id>T1490</id>
    </mitre>
  </rule>

  <rule id="100169" level="12">
    <if_sid>60612</if_sid>
    <field name="win.eventdata.param2">stopped</field>
    <field name="win.eventdata.param1">WinDefend|MsMpSvc</field>
    <description>Security service stopped</description>
    <mitre>
      <id>T1489</id>
    </mitre>
  </rule>

  <rule id="100170" level="10">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">del /f /s /q|rd /s /q</field>
    <description>Data destruction command</description>
    <mitre>
      <id>T1485</id>
    </mitre>
  </rule>

  <rule id="100171" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.commandLine">xmrig|stratum+tcp|cryptonight</field>
    <description>Cryptocurrency mining detected</description>
    <mitre>
      <id>T1496</id>
    </mitre>
  </rule>

  <rule id="100172" level="12">
    <if_group>sysmon_process_creation</if_group>
    <field name="sysmon.image">\\diskpart.exe</field>
    <field name="sysmon.commandLine">clean</field>
    <description>Disk wipe command detected</description>
    <mitre>
      <id>T1561</id>
    </mitre>
  </rule>

</group>
