# projects/legacy-archive/_archive/imports/soc-content/registry_run_keys.yml
title: Registry Run Key Persistence
id: 4d5e6f7a-8b9c-0d1e-2f3a-4b5c6d7e8f9a
status: stable
description: Detects modification of registry run keys commonly used for persistence
references:
    - https://attack.mitre.org/techniques/T1547/001/
author: HawkinsOps SOC
date: 2025/01/15
modified: 2025/01/15
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        TargetObject|contains:
            - '\Microsoft\Windows\CurrentVersion\Run\'
            - '\Microsoft\Windows\CurrentVersion\RunOnce\'
            - '\Microsoft\Windows\CurrentVersion\RunOnceEx\'
            - '\Microsoft\Windows\CurrentVersion\RunServices\'
            - '\Microsoft\Windows\CurrentVersion\Policies\Explorer\Run\'
    filter:
        Image|startswith:
            - 'C:\Program Files\'
            - 'C:\Program Files (x86)\'
        Details|contains:
            - 'Microsoft'
            - 'Windows Defender'
    condition: selection and not filter
falsepositives:
    - Legitimate software installations
    - System updates
level: high
