# projects/legacy-archive/_archive/imports/soc-content/scheduled_task_creation.yml
title: Suspicious Scheduled Task Creation
id: 3c4d5e6f-7a8b-9c0d-1e2f-3a4b5c6d7e8f
status: stable
description: Detects creation of scheduled tasks for persistence
references:
    - https://attack.mitre.org/techniques/T1053/005/
author: HawkinsOps SOC
date: 2025/01/15
modified: 2025/01/15
tags:
    - attack.persistence
    - attack.execution
    - attack.t1053.005
logsource:
    product: windows
    service: security
detection:
    selection:
        EventID: 4698
    filter_legitimage:
        TaskContent|contains:
            - '\Microsoft\Windows\'
            - 'MpCmdRun.exe'
    condition: selection and not filter_legitimage
falsepositives:
    - Legitimate software installations
    - Administrative tasks
    - System updates
level: medium
