# projects/legacy-archive/_archive/imports/soc-content/startup_folder_modification.yml
title: Startup Folder File Creation
id: 7a8b9c0d-1e2f-3a4b-5c6d-7e8f9a0b1c2d
status: stable
description: Detects file creation in Windows startup folders for persistence
references:
    - https://attack.mitre.org/techniques/T1547/001/
author: HawkinsOps SOC
date: 2025/01/15
modified: 2025/01/15
tags:
    - attack.persistence
    - attack.t1547.001
logsource:
    category: file_event
    product: windows
detection:
    selection:
        TargetFilename|contains:
            - '\Start Menu\Programs\Startup\'
            - '\Microsoft\Windows\Start Menu\Programs\StartUp\'
    condition: selection
falsepositives:
    - Legitimate software installations
    - User-created shortcuts
level: medium
