# projects/legacy-archive/_archive/imports/soc-content/detection-rules/sigma/persistence/linux_cron_persistence.yml
title: Linux Cron Job Persistence
id: 8b9c0d1e-2f3a-4b5c-6d7e-8f9a0b1c2d3e
status: stable
description: Detects suspicious cron job creation on Linux systems
references:
    - https://attack.mitre.org/techniques/T1053/003/
author: HawkinsOps SOC
date: 2025/01/15
modified: 2025/01/15
tags:
    - attack.persistence
    - attack.execution
    - attack.t1053.003
logsource:
    product: linux
    service: auditd
detection:
    selection_crontab:
        type: 'EXECVE'
        a0: 'crontab'
    selection_file:
        type: 'PATH'
        name|contains: '/etc/cron'
    condition: 1 of selection_*
falsepositives:
    - Legitimate scheduled tasks
    - System administration
    - Backup jobs
level: medium
