# projects/legacy-archive/_archive/imports/soc-content/detection-rules/sigma/persistence/linux_systemd_service.yml
title: Suspicious Systemd Service Creation
id: 9c0d1e2f-3a4b-5c6d-7e8f-9a0b1c2d3e4f
status: stable
description: Detects creation of suspicious systemd service units for persistence
references:
    - https://attack.mitre.org/techniques/T1543/002/
author: HawkinsOps SOC
date: 2025/01/15
modified: 2025/01/15
tags:
    - attack.persistence
    - attack.t1543.002
logsource:
    product: linux
    service: auditd
detection:
    selection:
        type: 'PATH'
        name|contains:
            - '/etc/systemd/system/'
            - '/usr/lib/systemd/system/'
            - '/.config/systemd/user/'
        name|endswith: '.service'
    condition: selection
falsepositives:
    - Legitimate service installations
    - Package manager operations
    - System administrators
level: medium
